# archive/starlinux/filesystem/etc/inittab
::sysinit:/etc/init.d/rcS
::once:cat /etc/welcome.txt
::respawn:/bin/cttyhack /bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:echo -e "\e[1;94m(****) \e[0mUnmounting All Filesystems..."
::shutdown:sync
::shutdown:umount -a -r
::shutdown:echo -e "\e[1;94m(****) \e[0mShutting Down StarLinux..."
::shutdown:sleep 1
::restart:/sbin/init
